# config/taproot.conf
description "Taproot Server"
author  "Chris Guilbeau <chris.guilbeau@gmail.com>"

start on runlevel [234]
stop on runlevel [0156]

chdir /root/taproot
exec /root/taproot/app.py
respawn
